# etc/udev/rules.d/60-keyboard-hotplug-setup.rules
ACTION=="add|bind", SUBSYSTEM=="usb|input", ATTR{idVendor}=="05ac", ATTR{idProduct}=="024f", RUN+="/home/stivius/scripts/system/setup_keyboard"
