# etc/udev/rules.d/60-monitor-hotplug-custom.rules
ACTION=="change", SUBSYSTEM=="drm", RUN+="update_monitors"
