# etc/udev/rules.d/60-monitor-hotplug-custom.rules
ACTION=="change", SUBSYSTEM=="drm", RUN+="/home/stivius/scripts/system/update_monitors"
