# etc/udev/rules.d/60-power-supply-notify.rules
ACTION=="change", SUBSYSTEM=="power_supply", ENV{POWER_SUPPLY_NAME}=="AC0", RUN+="check_battery_status"

